# php-debian/files/supervisor/conf.d/mailhog.conf
[program:mailhog]
command=/usr/local/bin/MailHog
directory=/usr/local/bin
user=root
autorestart=true
process_name=%(process_num)02d
numprocs=1